# deploy/backup/crontab
0 3 * * 0 /bin/sh /backup/backup.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
